# agent/systemd/chrome.service
[Unit]
Description=Google Chrome Kiosk
After=xvnc-chrome.service
Requires=xvnc-chrome.service

[Service]
Type=simple
Environment="DISPLAY=:1"
User=claworc
Environment="HOME=/home/claworc"
ExecStart=/usr/bin/google-chrome-stable --user-data-dir=/home/claworc/.config/google-chrome --remote-debugging-port=9222
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
