# agent/systemd/novnc-term.service
[Unit]
Description=noVNC Proxy for Terminal
After=xvnc-term.service
Requires=xvnc-term.service

[Service]
Type=simple
ExecStart=/opt/novnc/utils/novnc_proxy --vnc localhost:5902 --listen 6082
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
