# agent/systemd/novnc-chrome.service
[Unit]
Description=noVNC Proxy for Chrome
After=xvnc-chrome.service
Requires=xvnc-chrome.service

[Service]
Type=simple
User=claworc
ExecStart=/opt/novnc/utils/novnc_proxy --vnc localhost:5901 --listen 6081
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
