# agent/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway Service
After=network.target init-env.service
Wants=chrome.service

[Service]
Type=simple
Environment="DISPLAY=:1"
User=claworc
Environment="HOME=/home/claworc"
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
WorkingDirectory=/home/claworc
EnvironmentFile=-/etc/default/claworc
EnvironmentFile=-/etc/default/claworc-keys
ExecStartPre=/bin/bash -c 'test -n "$OPENCLAW_GATEWAY_TOKEN" && /usr/local/bin/openclaw config set gateway.auth.token "$OPENCLAW_GATEWAY_TOKEN" || true'
ExecStart=/usr/local/bin/openclaw gateway
Restart=always
RestartSec=3
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
