# agent/systemd/xvnc-chrome.service
[Unit]
Description=TigerVNC Server for Chrome (DISPLAY=:1)
After=network.target init-env.service
Requires=init-env.service

[Service]
Type=simple
User=claworc
EnvironmentFile=-/etc/default/vnc
ExecStart=/opt/scripts/start-xvnc.sh 1 5901
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
