# agent/rootfs/etc/s6-overlay/s6-rc.d/svc-sshd/run
#!/command/with-contenv bash

# Generate SSH host keys if they don't exist (first container start)
if [ ! -f /etc/ssh/ssh_host_ed25519_key ]; then
    echo "Generating SSH host keys..."
    ssh-keygen -A
fi

# Remove any legacy host key types that ssh-keygen -A may have created.
# We only keep Ed25519 (strongest) and RSA (widest compatibility).
rm -f /etc/ssh/ssh_host_dsa_key /etc/ssh/ssh_host_dsa_key.pub \
      /etc/ssh/ssh_host_ecdsa_key /etc/ssh/ssh_host_ecdsa_key.pub

# Create the privilege separation directory required by sshd
mkdir -p /run/sshd

# Ensure /root/.ssh directory exists with correct permissions
mkdir -p /root/.ssh
chmod 700 /root/.ssh

# Fix permissions on authorized_keys if it exists
if [ -f /root/.ssh/authorized_keys ]; then
    chmod 600 /root/.ssh/authorized_keys
fi

echo "Starting SSH daemon (hardened configuration)..."
# Run sshd in foreground mode (-D) so s6 can supervise it.
# -e sends logs to stderr instead of syslog; we redirect to the log file
# so the control plane's SSH-based log streamer can read them.
#
# SyslogFacility AUTH and LogLevel INFO are set in sshd_config.d/claworc.conf
# and control the log format/verbosity. To enable syslog-based centralized
# monitoring, remove the -e flag and ensure a syslog daemon (rsyslog, etc.)
# is running in the container.
exec /usr/sbin/sshd -D -e 2>> /var/log/sshd.log
