# agent/rootfs/etc/s6-overlay/s6-rc.d/svc-websockify/run
#!/bin/bash
exec websockify 0.0.0.0:5900 localhost:5901
