# agent/rootfs/etc/s6-overlay/s6-rc.d/svc-proxy/run
#!/command/with-contenv bash
# Write TLS cert and key from env vars to disk if provided
mkdir -p /config/ssl
if [ -n "$AGENT_TLS_CERT" ]; then
  printf '%s' "$AGENT_TLS_CERT" > /config/ssl/agent.crt
fi
if [ -n "$AGENT_TLS_KEY" ]; then
  printf '%s' "$AGENT_TLS_KEY" > /config/ssl/agent.key
  chmod 600 /config/ssl/agent.key
fi
if [ -n "$CONTROL_PLANE_CA" ]; then
  printf '%s' "$CONTROL_PLANE_CA" > /config/ssl/cp-ca.crt
fi
exec /usr/local/bin/claworc-proxy
