# agent/rootfs/etc/s6-overlay/s6-rc.d/svc-desktop/run
#!/bin/bash
export DISPLAY=:1 HOME=/home/claworc
exec s6-setuidgid claworc dbus-launch --exit-with-session startxfce4
